# peariso/pear/airootfs/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=/usr/bin/reflector --verbose --latest 5 --sort rate --save /etc/pacman.d/mirrorlist
ExecStart=-/sbin/agetty --autologin carly --noclear %I 38400 linux